# quadlet/api-product.container
[Unit]
Description=Product API Service (Exposed via SSL Proxy)
Documentation=https://github.com/your-org/api-product
After=internal-net.service
Requires=internal-net.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
Image=localhost/api-product:latest

ContainerName=api-product

# ============================================================================
# 網路配置
# ============================================================================
Network=internal-net

# Debug 模式端口映射
PublishPort=${DEBUG_PORT_8103}

# ============================================================================
# 環境變數
# ============================================================================
Environment=SERVICE_NAME=api-product
Environment=SERVICE_PORT=8080
Environment=ENVIRONMENT=production

# 資料庫連線
Environment=DB_HOST=your-db-host.example.com
Environment=DB_PORT=5432
Environment=DB_NAME=productdb

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=50%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=api-product
Label=tier=backend
Label=public-access=true
Label=exposed-via=ssl-proxy

[Service]
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
