# quadlet/api-product@.container
[Unit]
Description=Product API Service - Instance %i
After=internal-net.service postgres.service
Requires=internal-net.service

[Container]
Image=localhost/api-product:latest
ContainerName=api-product-%i
Network=internal-net:alias=api-product

Environment=SPRING_PROFILES_ACTIVE=prod

Volume=/opt/app/services/api-product/config/application-prod.yml:/app/config/application-prod.yml:ro
Secret=db-password,type=mount,target=/run/secrets/db-password.properties

HealthCmd=wget -qO- http://localhost:8080/actuator/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

Restart=always
RestartSec=10

Memory=768M
MemorySwap=768M
CPUQuota=50%

LogDriver=journald

Label=app=microservices
Label=service=api-product
Label=instance=%i
Label=tier=backend

[Service]
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
