# quadlet/api-user.container
[Unit]
Description=User API Service (Internal Only)
Documentation=https://github.com/your-org/api-user
After=internal-net.service
Requires=internal-net.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
# 替換成你的實際鏡像
# 範例：使用本地建置的鏡像
Image=localhost/api-user:latest

ContainerName=api-user

# ============================================================================
# 網路配置
# ============================================================================
# 只連接內部隔離網路
Network=internal-net

# Debug 模式端口映射（由環境變數控制）
# 生產環境：不設定此參數，完全隔離
# 開發環境：透過 .container.d/environment.conf 設定
PublishPort=${DEBUG_PORT_8101}

# ============================================================================
# 環境變數
# ============================================================================
Environment=SERVICE_NAME=api-user
Environment=SERVICE_PORT=8080
Environment=ENVIRONMENT=production

# 資料庫連線
Environment=DB_HOST=your-db-host.example.com
Environment=DB_PORT=5432
Environment=DB_NAME=userdb

# Java 應用程式設定（如適用）
# Environment=JAVA_OPTS=-Xms256m -Xmx512m -XX:+UseG1GC

# ============================================================================
# Volume 掛載
# ============================================================================
# Volume=/opt/app/api-user/config:/app/config:ro
# Volume=/opt/app/api-user/logs:/app/logs:rw

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=50%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=api-user
Label=tier=backend
Label=public-access=false

[Service]
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
