# quadlet/api-user@.container
[Unit]
Description=User API Service - Instance %i
Documentation=https://github.com/your-org/api-user
After=internal-net.service
Requires=internal-net.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
Image=localhost/api-user:latest

ContainerName=api-user-%i

# ============================================================================
# 網路配置
# ============================================================================
# 連接內部隔離網路，並設定 DNS 別名讓所有副本共用 "api-user" 名稱
# 其他服務（如 BFF）可直接透過 http://api-user:8080 訪問，
# Podman DNS 會自動 round-robin 到各副本
Network=internal-net:alias=api-user

# 多副本模式下不映射 debug port（避免衝突）
# 如需 debug 單一 instance，透過 podman exec 進入

# ============================================================================
# 環境變數
# ============================================================================
Environment=SERVICE_NAME=api-user
Environment=SERVICE_INSTANCE=%i
Environment=SERVICE_PORT=8080
Environment=ENVIRONMENT=production

# 資料庫連線由 environment.conf 提供（dev/prod 分開管理）

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=50%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=api-user
Label=instance=%i
Label=tier=backend
Label=public-access=false

[Service]
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
