# quadlet/bff@.container
[Unit]
Description=BFF Gateway (Spring Cloud Gateway) - Instance %i
After=internal-net.service
Requires=internal-net.service

[Container]
Image=localhost/bff:latest
ContainerName=bff-%i
Network=internal-net:alias=bff

Environment=SPRING_PROFILES_ACTIVE=prod

Volume=/opt/app/services/bff/config/application-prod.yml:/app/config/application-prod.yml:ro
Secret=jwt-signing-key,type=mount,target=/run/secrets/jwt-signing-key.properties

HealthCmd=curl -sf http://localhost:8080/actuator/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

Memory=768M
MemorySwap=768M
CPUQuota=50%

LogDriver=journald

Label=app=microservices
Label=service=bff
Label=instance=%i
Label=tier=gateway
Label=role=backend-for-frontend

[Service]
TimeoutStartSec=180
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
