# quadlet/bff@.container
[Unit]
Description=BFF (Backend For Frontend) Gateway - Instance %i
Documentation=https://github.com/your-org/bff-service
After=internal-net.service
Requires=internal-net.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
Image=localhost/bff:latest

ContainerName=bff-%i

# ============================================================================
# 網路配置
# ============================================================================
Network=internal-net:alias=bff

# ============================================================================
# 環境變數
# ============================================================================
Environment=SERVICE_NAME=bff
Environment=SERVICE_INSTANCE=%i
Environment=SERVICE_PORT=8080
Environment=ENVIRONMENT=production

# Backend API 端點（透過 DNS 別名，自動 round-robin 到各副本）
Environment=API_USER_URL=http://api-user:8080
Environment=API_ORDER_URL=http://api-order:8080
Environment=API_PRODUCT_URL=http://api-product:8080

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=50%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=bff
Label=instance=%i
Label=tier=gateway
Label=role=backend-for-frontend

[Service]
TimeoutStartSec=120
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
