# quadlet/ssl-proxy.container
[Unit]
Description=SSL Termination Proxy with Unified Token Verification
Documentation=https://github.com/your-org/microservices
After=internal-net.service frontend.service bff.service
Wants=frontend.service bff.service api-order.service api-product.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
# 使用 OpenResty（Nginx + LuaJIT）
Image=docker.io/openresty/openresty:alpine

ContainerName=ssl-proxy

# ============================================================================
# 網路配置
# ============================================================================
# 連接到 internal-net 以訪問 Backend APIs（如需直接訪問）
Network=internal-net

# 對外暴露 HTTP 和 HTTPS
PublishPort=80:80
PublishPort=443:443

# ============================================================================
# Volume 掛載
# ============================================================================
# 掛載自簽憑證（唯讀）
Volume=/opt/app/certs/server.crt:/etc/nginx/ssl/server.crt:ro
Volume=/opt/app/certs/server.key:/etc/nginx/ssl/server.key:ro

# 掛載 Nginx 配置
Volume=/opt/app/configs/ssl-proxy/nginx.conf:/usr/local/openresty/nginx/conf/nginx.conf:ro
Volume=/opt/app/configs/ssl-proxy/conf.d:/etc/nginx/conf.d:ro

# 掛載 Lua 模組
Volume=/opt/app/configs/ssl-proxy/lua:/usr/local/openresty/nginx/lua:ro

# 日誌目錄
Volume=/opt/app/logs/ssl-proxy:/var/log/nginx:rw

# ============================================================================
# Partner Secrets 配置（混合方案）
# ============================================================================
# 生產環境：使用 Podman Secrets（優先）
# 如果 Secret 存在，會注入為環境變數並覆蓋下方的 Environment
Secret=jwt-secret-partner-a,type=env,target=JWT_SECRET_PARTNER_A
Secret=jwt-secret-partner-b,type=env,target=JWT_SECRET_PARTNER_B
Secret=jwt-secret-partner-c,type=env,target=JWT_SECRET_PARTNER_C

# 開發環境：Fallback 環境變數（Secret 不存在時使用）
# 這些會被 .container.d/environment.conf 覆蓋
Environment=JWT_SECRET_PARTNER_A=dev-secret-partner-a-fallback-32chars
Environment=JWT_SECRET_PARTNER_B=dev-secret-partner-b-fallback-32chars
Environment=JWT_SECRET_PARTNER_C=dev-secret-partner-c-fallback-32chars

# Debug 模式（開發環境可透過 .container.d 覆蓋）
Environment=DEBUG_PARTNERS=false

# Nginx 設定
Environment=NGINX_WORKER_PROCESSES=auto
Environment=NGINX_WORKER_CONNECTIONS=2048

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -fk https://localhost/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=15s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=100%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=ssl-proxy
Label=tier=gateway
Label=role=ssl-termination

[Service]
TimeoutStartSec=60
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
