# quadlet/ssl-proxy.container
[Unit]
Description=SSL Termination Proxy with Unified Token Verification
Documentation=https://github.com/your-org/microservices
After=internal-net.service frontend.service bff.service
Wants=frontend.service bff.service api-order.service api-product.service

[Container]
# ============================================================================
# 鏡像配置
# ============================================================================
# 使用 OpenResty（Nginx + LuaJIT）
Image=docker.io/openresty/openresty:alpine

ContainerName=ssl-proxy

# ============================================================================
# 網路配置
# ============================================================================
# 連接到 internal-net 以訪問 Backend APIs（如需直接訪問）
Network=internal-net

# 對外暴露 HTTP 和 HTTPS
PublishPort=80:80
PublishPort=443:443

# ============================================================================
# Volume 掛載
# ============================================================================
# 掛載自簽憑證（唯讀）
Volume=/opt/app/certs/server.crt:/etc/nginx/ssl/server.crt:ro
Volume=/opt/app/certs/server.key:/etc/nginx/ssl/server.key:ro

# 掛載 Nginx 配置
Volume=/opt/app/configs/ssl-proxy/nginx.conf:/usr/local/openresty/nginx/conf/nginx.conf:ro
Volume=/opt/app/configs/ssl-proxy/conf.d:/etc/nginx/conf.d:ro

# 掛載 Lua 腳本
Volume=/opt/app/configs/ssl-proxy/lua:/usr/local/openresty/nginx/lua:ro

# 日誌目錄
Volume=/opt/app/logs/ssl-proxy:/var/log/nginx:rw

# ============================================================================
# 環境變數
# ============================================================================
# JWT Secret（請更換成強隨機字串）
Environment=JWT_SECRET=your-super-secret-key-change-me-in-production-at-least-32-chars

# Nginx 設定
Environment=NGINX_WORKER_PROCESSES=auto
Environment=NGINX_WORKER_CONNECTIONS=2048

# ============================================================================
# 健康檢查
# ============================================================================
HealthCmd=curl -fk https://localhost/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=15s

# ============================================================================
# 重啟策略
# ============================================================================
Restart=always
RestartSec=10

# ============================================================================
# 資源限制
# ============================================================================
Memory=512M
MemorySwap=512M
CPUQuota=100%

# ============================================================================
# 日誌配置
# ============================================================================
LogDriver=journald

# ============================================================================
# 標籤
# ============================================================================
Label=app=microservices
Label=service=ssl-proxy
Label=tier=gateway
Label=role=ssl-termination

[Service]
TimeoutStartSec=60
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target default.target
